# daemons/autopilot-tick.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Local Autopilot — macOS launchd template.

  TO INSTALL:
    1. Copy this file to ~/Library/LaunchAgents/local-autopilot-tick.plist
    2. Replace every `__REPO_DIR__` placeholder with the absolute path to your
       Local Autopilot clone (e.g. /Users/aarontjomsland/dev/local-autopilot)
    3. Replace `__USER__` with your macOS username (output of `whoami`)
    4. Optional: tune `StartInterval` (seconds between ticks; default 1800 = 30 min)
    5. Load it:
         launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/local-autopilot-tick.plist
    6. Verify:
         launchctl list | grep local-autopilot
         tail -f ~/.context-dna/autopilot-logs/launchd.out

  TO UNINSTALL:
    launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/local-autopilot-tick.plist
    rm ~/Library/LaunchAgents/local-autopilot-tick.plist

  The tick runs `archloop_runner.py --cycles 1` every interval. The runner
  itself respects the autopilot state machine — if mode == "off" it exits
  immediately. Safe to leave loaded permanently; toggle activity via
  `autopilot on` / `autopilot off`.

  ──────────────────────────────────────────────────────────────────────────
  OPTIONAL: HEADLESS EXECUTOR (opt-in via installer, NOT baked into this template)
  ──────────────────────────────────────────────────────────────────────────
  By default, the runner emits agent prompts and stops at AWAIT — you (the
  human) read the proposed action and decide whether to run it. This is the
  safe default and is what this template renders.

  If you add `--headless-executor` to the ProgramArguments array, Atlas will
  invoke `claude --print` on every cycle prompt and execute the agent
  without your review. That trades safety for full autonomy:
    + bigger trust delegation (Atlas decides AND does)
    + faster turnaround (no human-in-the-loop wait)
    - no human checkpoint per cycle
    - cost + side-effects bounded only by the per-agent budget / timeout /
      `--bare` rails inside the runner
    - cycle artifacts under ~/.context-dna/autopilot-logs/cycle-* are now
      your sole audit trail — REVIEW THEM REGULARLY.

  DO NOT edit this template to add the flag. Instead, opt in via the
  installer so the choice is visible in shell history and reversible:
      bash scripts/install-daemon-macos.sh --headless-executor
  Revert by re-running the installer without the flag, or `autopilot off`.

  ──────────────────────────────────────────────────────────────────────────
  OPTIONAL: FLEET_ESCALATE (opt-in via installer, NOT baked into this template)
  ──────────────────────────────────────────────────────────────────────────
  When all `--cycles` finish without `satisfied=True`, the runner can POST
  the unresolved cycle to the Multi-Fleet chief node (mac1) at
  http://127.0.0.1:8855/message so a different vendor LLM with deeper
  ContextDNA can break the tie or return UNRESOLVABLE.

  Opt in via the installer (do NOT edit this template):
      bash scripts/install-daemon-macos.sh --fleet-escalate
  Combine with --headless-executor for fully-autonomous chief-arbitrated runs.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.localautopilot.tick</string>

  <key>ProgramArguments</key>
  <array>
    <string>__REPO_DIR__/.venv/bin/python3</string>
    <string>-m</string>
    <string>local_autopilot.tools.archloop_runner</string>
    <string>--cycles</string>
    <string>1</string>
    <string>--cost-cap-usd</string>
    <string>5.00</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_DIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
    <key>PYTHONPATH</key>
    <string>__REPO_DIR__:__REPO_DIR__/local_autopilot</string>
  </dict>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/__USER__/.context-dna/autopilot-logs/launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/__USER__/.context-dna/autopilot-logs/launchd.err</string>

  <!-- Be a good citizen: low CPU/IO/memory priority so we never starve
       interactive work. The runner is bursty but should never spike. -->
  <key>ProcessType</key>
  <string>Background</string>
  <key>Nice</key>
  <integer>10</integer>
</dict>
</plist>
